# deploy/skipper.service
[Unit]
Description=Skipper Agent OS Daemon
Documentation=https://github.com/sethdford/skipper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=skipper
Group=skipper
ExecStart=/usr/local/bin/skipper start
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Environment
EnvironmentFile=-/etc/skipper/env
WorkingDirectory=/var/lib/skipper

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/skipper
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=false
RestrictRealtime=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
